# cron/daily_job.cron
58 1 * * * root cd /app && export DJANGO_SETTINGS_MODULE=DesafioVeeries.settings && /usr/local/bin/python manage.py bronze_stage_paranagua && /usr/local/bin/python manage.py bronze_stage_santos && /usr/local/bin/python manage.py silver_stage >> /var/log/daily_job.log 2>&1
